# crontab
* * * * * sudo python home/report.py >> temperature.log
