# crontab
* * * * * sudo python report.py >> temperature.log
